# rrd.service
[Unit]
Description=Run rrdtool

[Service]
Type=simple
ExecStart=/home/alarm/rrd.sh
